# freenom_updater.service
[Unit]

Description=freenom_updater

After=syslog.target

[Serive]

ExecStart=/opt/freenom_updater.sh

[Install]

WantedBy=multi-user.target
